# base/etc/systemd/system/tangd.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=51379 